# systemd/user/hyprpaper.service
[Unit]
Description=Hyprpaper wallpaper daemon
After=graphical-session.target

[Service]
ExecStart=/usr/bin/hyprpaper
Restart=on-failure

[Install]
WantedBy=graphical-session.target



